# systemd/coc-report-deploy.service
[Unit]
Description=CoC Report Static Site Deploy
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/clashcommand/app
EnvironmentFile=/opt/clashcommand/app/.env
ExecStart=/opt/clashcommand/app/scripts/deploy_report_once.sh
