# systemd/coc-report-deploy.service
[Unit]
Description=DEPRECATED duplicate CoC Report Static Site Deploy
Documentation=file:/opt/clashcommand/app/systemd/README.md
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/clashcommand/app
EnvironmentFile=/opt/clashcommand/app/.env
ExecStart=/opt/clashcommand/app/scripts/deploy_report_once.sh
